# configs/etc/supervisor.d/blog.ini
[program:blog]
command=/srv/services/blog/blog
stderr_logfile = /var/log/supervisor/blog-stderr.log
stdout_logfile = /var/log/supervisor/blog-stdout.log
